# autostart/runStartupHacks.desktop
[Desktop Entry]
Type=Application
Name=run-startup-hacks
Exec=/home/user1/startupHacks.sh

